# config/supervisor/supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:fpm]
command=/usr/sbin/php7-fpm -R --nodaemonize

[program:ssh]
command=/usr/sbin/sshd -D